# stockbot.service
[Unit]
Description=Stock Telegram Bot (Gemini AI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=seokhwanlee3
WorkingDirectory=/home/seokhwanlee3/stock-bot
ExecStart=/home/seokhwanlee3/stock-bot/venv/bin/python telegram_bot.py
EnvironmentFile=/home/seokhwanlee3/stock-bot/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
