# stockbot.service
[Unit]
Description=Stock Telegram Bot (Gemini AI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/stock-bot
ExecStart=/usr/bin/python3 telegram_bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
